# devel/ansible/roles/dev/files/pw-import.path
[Unit]
Description=Loads mail into Patchwork

[Path]
PathChanged=/var/spool/mail/vagrant

[Install]
WantedBy=default.target
