# devel/ansible/roles/patchlab/files/certbot-renew.timer
[Unit]
Description=Renew Certbot certificate

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1200

[Install]
WantedBy=multi-user.target
