# devel/ansible/roles/patchlab/files/certbot-renew.service
[Unit]
Description=Renew Certbot certificate
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/certbot renew
